# config/devicehub-provisioning.service
[Unit]
Description=Edgeberry Device Hub - Provisioning Service
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/Edgeberry/devicehub/provisioning-service
EnvironmentFile=-/etc/Edgeberry/devicehub/provisioning.env
Environment=NODE_ENV=production
Environment=MQTT_URL=mqtt://127.0.0.1:1883
Environment=DEVICEHUB_DB=/var/lib/edgeberry/devicehub/devicehub.db
ExecStart=/usr/bin/node /opt/Edgeberry/devicehub/provisioning-service/dist/index.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
