# config/devicehub-provisioning.service
[Unit]
Description=Edgeberry Device Hub - Provisioning Service
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/Edgeberry/devicehub/provisioning-service
Environment=NODE_ENV=production
Environment=MQTT_URL=mqtts://localhost:8883
Environment=MQTT_TLS_CA=/opt/Edgeberry/devicehub/config/certs/ca.crt
Environment=MQTT_TLS_CERT=/opt/Edgeberry/devicehub/config/certs/provisioning.crt
Environment=MQTT_TLS_KEY=/opt/Edgeberry/devicehub/config/certs/provisioning.key
Environment=MQTT_TLS_REJECT_UNAUTHORIZED=true
EnvironmentFile=-/etc/Edgeberry/devicehub/provisioning.env
ExecStart=/usr/bin/node /opt/Edgeberry/devicehub/provisioning-service/dist/index.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
