# config/devicehub-registry.service
[Unit]
Description=Edgeberry Device Hub - Device Registry Service
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/Edgeberry/devicehub/registry-service
Environment=NODE_ENV=production
Environment=MQTT_URL=mqtts://localhost:8883
Environment=MQTT_TLS_CA=/opt/Edgeberry/devicehub/config/certs/ca.crt
Environment=MQTT_TLS_CERT=/opt/Edgeberry/devicehub/config/certs/registry.crt
Environment=MQTT_TLS_KEY=/opt/Edgeberry/devicehub/config/certs/registry.key
Environment=MQTT_TLS_REJECT_UNAUTHORIZED=true
EnvironmentFile=-/etc/Edgeberry/devicehub/registry.env
ExecStart=/usr/bin/node /opt/Edgeberry/devicehub/registry-service/dist/index.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
